# deploy/systemd/chatbot-watchdog.timer
[Unit]
Description=Chatbot health check every 90s (fixes offline ngrok after suspend)

[Timer]
OnBootSec=2min
OnUnitActiveSec=90s
AccuracySec=30s
Persistent=true

[Install]
WantedBy=timers.target
