# deploy/systemd/chatbot-watchdog.service
[Unit]
Description=Chatbot health check (backend + ngrok)
After=network-online.target docker.service

[Service]
Type=oneshot
ExecStart=/home/jadaboawwad/Files/Software/Repositories/Applications/Chatbot/deploy/systemd/chatbot-watchdog.sh
User=jadaboawwad
Group=jadaboawwad
