# deploy/systemd/chatbot-ngrok.service
[Unit]
Description=Chatbot - ngrok tunnel to gateway (port 9080)
After=network-online.target chatbot.service
Wants=network-online.target
Requires=chatbot.service

[Service]
Type=simple
ExecStart=/usr/local/bin/ngrok start chatbot
User=jadaboawwad
Group=jadaboawwad
StandardOutput=journal
StandardError=journal
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
