# deploy/ngrok/chatbot-ngrok.service
[Unit]
Description=ngrok tunnel to chatbot gateway (port 9080)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ngrok start chatbot
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
