# deploy/ngrok/chatbot-ngrok.service
[Unit]
Description=ngrok tunnel to Streamlit (port 8501)
After=network-online.target docker.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/ngrok http 8501 --host-header=rewrite
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
